# s6/web/run
#!/bin/sh

exec 2>&1
echo "web started"
exec '/web/server'
